# systemd-user/.config/systemd/user/LG_Buddy_screen.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/local/bin/LG_Buddy_Screen_Monitor_GNOME
